# packaging/app-software-repository.cron
0 5 * * * root /usr/clearos/apps/software_repository/deploy/config-check
